# clipboardnotificator.service
[Unit]
Description = Clipboard notificator

[Service]
ExecStart = commad_path
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target
